# cron/crontab
#crontab
* * * * * root *applications/OpenLex_pruebs/cron/check_email.py
