# cron/crontab
#crontab
0 0 * * * root *applications/OpenLex_pruebs/cron/check_email.py
